# deploy/systemd/system/opencloud-server.service
[Unit]
Description=OpenCloud API server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__OPENCLOUD_SERVICE_USER__
EnvironmentFile=__OPENCLOUD_SERVICE_ENV__
ExecStart=/bin/bash -lc 'cd "$$OPENCLOUD_REPO_DIR" && "$${PNPM_BIN:-pnpm}" run start --filter=server'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opencloud-server

[Install]
WantedBy=multi-user.target
