# deploy/systemd/user/opencloud-nova.service
[Unit]
Description=OpenCloud Nova web UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/opencloud/opencloud-service.env
ExecStart=/bin/bash -lc 'cd "$OPENCLOUD_REPO_DIR" && pnpm run start --filter=nova'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opencloud-nova

[Install]
WantedBy=default.target
